# dist/http-tarpit.service
[Unit]
Description=HTTP tarpit
Wants=network-online.target
After=network.target network-online.target

[Service]
User=root
Type=simple
WorkingDirectory=/tmp/
PrivateTmp=true
# bind to default http port, expose metrics port for a trusted local network
ExecStart=/usr/local/bin/http_tarpit -ssh=":80" -prom="10.0.0.123:5002"
RestartSec=5

# keep retrying no matter what
Restart=always

[Install]
WantedBy=multi-user.target